# weather-station.service
[Unit]
Description=Raspberry Pi Weather Station with E-Ink Display
Documentation=https://github.com/SpeedySwifter/waveshare2in13b_V4---WeatherStation
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=svenhajer
Group=svenhajer
WorkingDirectory=/Users/svenhajer/Documents/GitHub/waveshare2in13b_V4---WeatherStation
ExecStart=/usr/bin/python3 /Users/svenhajer/Documents/GitHub/waveshare2in13b_V4---WeatherStation/weather_station.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration - very robust
Restart=always
RestartSec=15
TimeoutStartSec=60
TimeoutStopSec=30

# Failure handling
StartLimitAction=reboot
FailureAction=restart

# Process management
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-station

# Security and resource limits
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/Users/svenhajer/Documents/GitHub/waveshare2in13b_V4---WeatherStation
MemoryMax=256M
CPUQuota=50%

# Environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
